# sleep-app.service
[Unit]
Description=Sleep main daemon

[Service]
ExecStart=/usr/local/bin/sleep-app

[Install]
WantedBy=multi-user.target
